# noisebell.service
[Unit]
Description=Noisebell Discord Notification Service
After=network.target

[Service]
Type=simple
User=noisebridge
WorkingDirectory=/home/noisebridge/Documents/noisebell
Environment=DISCORD_TOKEN=your_bot_token_here
Environment=DISCORD_CHANNEL_ID=your_channel_id_here
ExecStart=/home/noisebridge/.cargo/bin/noisebell
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 